# DockTile/Resources/LoginItems/com.docktile.app.tilelauncher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<!-- Release launcher agent. Registered via SMAppService.agent(plistName:). -->
	<!-- Runs the main app binary headless at login to warm the visible tiles. -->
	<key>Label</key>
	<string>com.docktile.app.tilelauncher</string>
	<key>BundleProgram</key>
	<string>Contents/MacOS/Dock Tile</string>
	<key>ProgramArguments</key>
	<array>
		<string>Contents/MacOS/Dock Tile</string>
		<string>--login-spawn-tiles</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>LimitLoadToSessionType</key>
	<string>Aqua</string>
</dict>
</plist>
